# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid  ; 将 PID 文件放置在 /var/run/ 目录
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock


[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/var/log/nginx/nginx_stdout.log
stderr_logfile=/var/log/nginx/nginx_stderr.log
autostart=true
autorestart=true

[program:webman]
command=php /var/www/start.php start
stdout_logfile=/var/log/webman/webman_stdout.log
stderr_logfile=/var/log/webman/webman_stderr.log
autostart=true
autorestart=true